# crontab
* * * * * root python3 py_twitter-posts-scraper/twitter-post-scraper.py >> /logfile.txt
