# crontab
*/2 * * * * python3 twitter-post-scraper.py >> logfile.txt
